# deploy/systemd/content-ops-runtime-backup.timer.example
[Unit]
Description=Run Content Ops runtime backup daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=content-ops-runtime-backup.service

[Install]
WantedBy=timers.target
